# deploy/construct-ice.service
[Unit]
Description=Construct ICE Relay (construct-ice)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=construct
Group=construct
WorkingDirectory=/opt/construct-ice
ExecStart=/usr/local/bin/construct-ice
EnvironmentFile=/etc/construct-ice/env
Restart=always
RestartSec=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=construct-ice

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/construct-ice
PrivateTmp=true

[Install]
WantedBy=multi-user.target
